# backup/backup@data.timer
[Unit]
Description=Backs up data daily

[Timer]
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target