# backup/backupserver@globaltoaster.timer
[Unit]
Description=Backs up server daily

[Timer]
OnCalendar=daily
Persistent=true
 
[Install]
WantedBy=timers.target